# s6-overlay/s6-rc.d/05-snapserver/run
#!/command/with-contenv sh
echo "Starting snapserver"
exec /usr/local/bin/snapserver
